# files/docker-keycloak-postgres-pgdump.timer
[Unit]
Description=Daily backup of Postgres database

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
